# docker/files/wlp.ini
[program:wlp]
user=root
command=/opt/websphere/wlp/bin/server start
autorestart=false