# docker/files/awslogs/awslog.ini
[program:awslogs]
command=/usr/bin/env -i AWS_CONFIG_FILE=/etc/aws.conf /usr/bin/aws logs push --config-file /etc/awslogs-agent.conf
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null