# docker/files/awslogs/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
numprocs=1
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog